# stanford-stroke-pacs/launchd/com.ssc.colima.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Boot persistence + crash recovery for the Colima Docker VM (headless macOS).

  Runs scripts/macos/colima_watchdog.sh, a long-lived supervisor that brings the
  VM up at boot (via the idempotent colima_start.sh) and restarts it within
  CHECK_INTERVAL seconds if it ever crashes or stops. Orthanc (restart:
  unless-stopped) then returns on its own once dockerd is back.

  KeepAlive is TRUE here because the watchdog is a never-exiting process: if the
  watchdog itself is killed, launchd revives it after ThrottleInterval. (Do NOT
  point this at colima_start.sh directly — that exits 0 once the VM is up, which
  KeepAlive would treat as "finished" and relaunch in a busy loop.)

  Install as a LaunchDaemon (system domain — required, this box is headless with
  no login session): sudo scripts/macos/install_launchd.sh
-->
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.ssc.colima</string>
  <key>UserName</key>         <string>pere</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/ssc-pacs/ssc-pacs/stanford-stroke-pacs/scripts/macos/colima_watchdog.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>/Users/pere</string>
  </dict>
  <key>RunAtLoad</key>        <true/>
  <key>KeepAlive</key>        <true/>
  <key>ThrottleInterval</key> <integer>30</integer>
  <key>StandardOutPath</key>  <string>/Users/pere/Library/Logs/ssc-colima.log</string>
  <key>StandardErrorPath</key><string>/Users/pere/Library/Logs/ssc-colima.err</string>
</dict>
</plist>
