# stanford-stroke-pacs/launchd/com.ssc.pg-backup-freshness.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ssc.pg-backup-freshness</string>
  <key>UserName</key><string>pere</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/opt/coreutils/libexec/gnubin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>/Users/pere</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>LC_ALL</key><string>en_US.UTF-8</string>
  </dict>
  <key>StandardOutPath</key><string>/Users/pere/Library/Logs/com.ssc.pg-backup-freshness.log</string>
  <key>StandardErrorPath</key><string>/Users/pere/Library/Logs/com.ssc.pg-backup-freshness.err</string>
  <key>ProgramArguments</key><array><string>/opt/homebrew/bin/bash</string><string>/opt/ssc-pacs/ssc-pacs/stanford-stroke-pacs/scripts/backup/check_backup_freshness.sh</string></array>
  <key>StartCalendarInterval</key><dict><key>Minute</key><integer>0</integer></dict>
  <key>RunAtLoad</key><true/>
</dict></plist>
