# stanford-stroke-pacs/launchd/com.ssc.webapp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Web App (uvicorn on :8043) service — macOS replacement for systemd/ssc-web-app.service.

  HEADLESS NOTE: this Mac has no console/GUI login session, so a per-user
  LaunchAgent (gui/$UID domain) will NOT load (`brew services`/`launchctl
  bootstrap gui/...` fail with "Domain does not support specified action").
  Install this as a LaunchDaemon (system domain, loaded at boot) instead:

    sudo cp launchd/com.ssc.webapp.plist /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/com.ssc.webapp.plist
    sudo launchctl bootstrap system /Library/LaunchDaemons/com.ssc.webapp.plist
    sudo launchctl kickstart -k system/com.ssc.webapp     # restart
    sudo launchctl print system/com.ssc.webapp            # status

  It runs as user `pere` (UserName) so it shares Postgres loopback-trust and the
  conda env. Depends on Postgres + Colima being up first (com.ssc.colima +
  Postgres LaunchDaemon) — KeepAlive restarts it until the DB pool connects.
-->
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.ssc.webapp</string>
  <key>UserName</key>         <string>pere</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/Caskroom/miniconda/base/envs/ssc-pacs/bin/uvicorn</string>
    <string>app:app</string>
    <string>--host</string><string>0.0.0.0</string>
    <string>--port</string><string>8043</string>
  </array>
  <key>WorkingDirectory</key> <string>/opt/ssc-pacs/ssc-pacs/stanford-stroke-pacs/web-app</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>/Users/pere</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>LC_ALL</key><string>en_US.UTF-8</string>
  </dict>
  <key>RunAtLoad</key>        <true/>
  <key>KeepAlive</key>        <true/>
  <key>StandardOutPath</key>  <string>/Users/pere/Library/Logs/ssc-web-app.log</string>
  <key>StandardErrorPath</key><string>/Users/pere/Library/Logs/ssc-web-app.err</string>
</dict>
</plist>
