# stanford-stroke-pacs/launchd/com.ssc.cold-storage-health.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ssc.cold-storage-health</string>
  <key>UserName</key><string>__DEPLOY_USER__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__BREW_PREFIX__/opt/coreutils/libexec/gnubin:__BREW_PREFIX__/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>__HOME_DIR__</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>LC_ALL</key><string>en_US.UTF-8</string>
  </dict>
  <key>StandardOutPath</key><string>__HOME_DIR__/Library/Logs/com.ssc.cold-storage-health.log</string>
  <key>StandardErrorPath</key><string>__HOME_DIR__/Library/Logs/com.ssc.cold-storage-health.err</string>
  <key>ProgramArguments</key><array><string>__CONDA_ENV_BIN__/python</string><string>__REPO_ROOT__/scripts/cold_storage/cold_storage_health.py</string><string>--quiet</string></array>
  <key>StartInterval</key><integer>900</integer>
  <key>RunAtLoad</key><false/>
</dict></plist>
