# stanford-stroke-pacs/deploy/launchd/com.ssc.orthanc-storage-backup.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.ssc.orthanc-storage-backup</string>
  <key>UserName</key><string>__DEPLOY_USER__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>__BREW_PREFIX__/opt/coreutils/libexec/gnubin:__BREW_PREFIX__/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key><string>__HOME_DIR__</string>
    <key>LANG</key><string>en_US.UTF-8</string>
    <key>LC_ALL</key><string>en_US.UTF-8</string>
    <key>DOCKER_HOST</key><string>unix://__HOME_DIR__/.colima/default/docker.sock</string>
  </dict>
  <key>StandardOutPath</key><string>__HOME_DIR__/Library/Logs/com.ssc.orthanc-storage-backup.log</string>
  <key>StandardErrorPath</key><string>__HOME_DIR__/Library/Logs/com.ssc.orthanc-storage-backup.err</string>
  <key>ProgramArguments</key><array><string>__BREW_PREFIX__/bin/bash</string><string>__REPO_ROOT__/scripts/backup/backup_orthanc_storage.sh</string></array>
  <key>StartCalendarInterval</key><dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>45</integer></dict>
  <key>RunAtLoad</key><false/>
</dict></plist>
