# stanford-stroke-pacs/systemd/orthanc-storage-backup.timer
[Unit]
Description=PACS — nightly backup of the Orthanc storage volume
Documentation=file:/home/perecanals/ssc-pacs/stanford-stroke-pacs/documentation/operations/backup_strategy.md

[Timer]
# Daily at 02:45 local time — after the PostgreSQL dumps (02:15 / 02:30) so the
# three backups don't collide on disk I/O. Persistent=true catches up a missed
# run (host off) on next boot.
OnCalendar=*-*-* 02:45:00
RandomizedDelaySec=10m
Persistent=true
Unit=orthanc-storage-backup.service

[Install]
WantedBy=timers.target
