# stanford-stroke-pacs/systemd/pg-backup-stanford-stroke.service
[Unit]
Description=PACS — pg_dump backup of stanford-stroke DB
Documentation=file:/home/perecanals/pacs/stanford-stroke-pacs/documentation/operations/backup_strategy.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=perecanals
Group=perecanals
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ExecStart=/home/perecanals/pacs/stanford-stroke-pacs/scripts/backup_pg_db.sh stanford-stroke
# Restart-on-failure has no meaning for oneshot timer-driven units; the
# next scheduled run handles transient failures and OnFailure= alerts in
# the meantime.
StandardOutput=journal
StandardError=journal
