# stanford-stroke-pacs/systemd/reconciliation.timer.in
[Unit]
Description=PACS — two-DB reconciliation every 6 hours
Documentation=file:__REPO_ROOT__/documentation/operations/reconciliation.md

[Timer]
OnBootSec=10min
OnUnitActiveSec=6h
RandomizedDelaySec=5min
Persistent=true
Unit=reconciliation.service

[Install]
WantedBy=timers.target
