# stanford-stroke-pacs/systemd/reconciliation.service.in
[Unit]
Description=PACS — two-DB reconciliation (image_series vs Orthanc index)
Documentation=file:__REPO_ROOT__/documentation/operations/reconciliation.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__DEPLOY_USER__
Group=__DEPLOY_GROUP__
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ExecStart=__PYTHON_BIN__ __REPO_ROOT__/scripts/data_integrity/reconcile.py --json --quiet
StandardOutput=journal
StandardError=journal
