# stanford-stroke-pacs/systemd/ssc-web-app.service.in
[Unit]
Description=SSC Web App (FastAPI)
After=network.target postgresql.service

[Service]
Type=simple
User=__DEPLOY_USER__
WorkingDirectory=__REPO_ROOT__/web-app
ExecStart=__UVICORN_BIN__ app:app --host 0.0.0.0 --port 8043
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
