# stanford-stroke-pacs/deploy/systemd/cold-storage-health.timer.in
[Unit]
Description=PACS — cold-storage health probe every 15 minutes
Documentation=file:__DOCS_ROOT__/cold_storage/runbook.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
RandomizedDelaySec=2min
Persistent=true
Unit=cold-storage-health.service

[Install]
WantedBy=timers.target
