# stanford-stroke-pacs/deploy/systemd/pg-backup-orthanc.service.in
[Unit]
Description=PACS — pg_dump backup of orthanc_db
Documentation=file:__DOCS_ROOT__/operations/backup_strategy.md
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=__DEPLOY_USER__
Group=__DEPLOY_GROUP__
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ExecStart=__REPO_ROOT__/scripts/backup/backup_pg_db.sh orthanc_db
StandardOutput=journal
StandardError=journal
