# stanford-stroke-pacs/deploy/systemd/pg-backup-freshness.timer.in
[Unit]
Description=PACS — hourly backup freshness check
Documentation=file:__DOCS_ROOT__/operations/backup_strategy.md

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true
Unit=pg-backup-freshness.service

[Install]
WantedBy=timers.target
